# docker/web/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D
redirect_stderr = true
autorestart = true

[program:apache2]
command=/usr/local/bin/apache2-foreground
redirect_stderr = true
autorestart = true